# root/etc/services.d/updater/run
#!/usr/bin/with-contenv bash

exec /app/yt-dlp/updater.sh